# app/app-stop.conf
description "Mognet server stop"

start on runlevel [06]

script
  cd /var/www/mognet/
  exec /usr/bin/npm stop
end script